# deploy/com.munibot.ollama.plist.template
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
  <dict>
    <key>Label</key>
    <string>com.munibot.ollama</string>
    <key>ProgramArguments</key>
    <array>
      <string>/usr/local/bin/ollama</string>
      <string>serve</string>
    </array>
    <key>EnvironmentVariables</key>
    <dict>
      <key>OLLAMA_HOST</key>
      <string>0.0.0.0:11434</string>
    </dict>
    <key>RunAtLoad</key>
    <true/>
    <key>KeepAlive</key>
    <true/>
    <key>StandardOutPath</key>
    <string>__HOME__/Library/Logs/munibot-ollama.log</string>
    <key>StandardErrorPath</key>
    <string>__HOME__/Library/Logs/munibot-ollama.err.log</string>
    <key>WorkingDirectory</key>
    <string>__HOME__</string>
  </dict>
</plist>
